# debian/hc05-configurator.service
[Unit]
Description=HC-05 Configurator GUI

[Service]
Type=simple
ExecStart=/usr/bin/hc05-configurator

[Install]
WantedBy=graphical.target
